# services/sshd-rendong/run
#!/bin/sh
exec 2>&1

HOME="/home/username"

exec /usr/sbin/sshd -D \
    -f "$HOME/.ssh/custom.conf"
